# wan-watcher.service.j2
[Unit]
Description=WAN-Watcher
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash {{ ww_script_dir }}/wan-watcher.sh
